# var/lib/dpkg/info/mariadb-server.conffiles
/etc/apparmor.d/usr.sbin.mariadbd
/etc/init.d/mariadb
/etc/logcheck/ignore.d.paranoid/mariadb-server
/etc/logcheck/ignore.d.server/mariadb-server
/etc/logcheck/ignore.d.workstation/mariadb-server
/etc/logrotate.d/mariadb
/etc/mysql/debian-start
/etc/mysql/mariadb.conf.d/50-mysqld_safe.cnf
/etc/mysql/mariadb.conf.d/50-server.cnf
/etc/security/user_map.conf
